# phosphor-multi-gpio-monitor@.service
[Unit]
Description=Phosphor %I GPIO monitor

[Service]
Restart=no
ExecStart=/usr/bin/phosphor-multi-gpio-monitor --config {sysconfdir}/%I.json
